# internal/crewctl/testdata/launchd_promo.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.shipyard.crew.promo</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/shipyard-crew</string>
    <string>--service</string>
    <string>--agent</string>
    <string>promo</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>HOME/.shipyard/logs/crew/promo.out.log</string>
  <key>StandardErrorPath</key><string>HOME/.shipyard/logs/crew/promo.err.log</string>
</dict>
</plist>
